# tests/images/redis/supervisor_redis.conf
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
process_name=%(program_name)s
autostart=true
autorestart=false
stopsignal=TERM
priority=5
redirect_stderr=true
stdout_logfile=/var/log/redis/server.log
