# tests/images/redis/supervisor_senticache.conf
[program:senticache]
command=/usr/bin/redis-senticache /etc/redis/senticache.conf
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=TERM
priority=5
redirect_stderr=true
stdout_logfile=/var/log/redis/senticache.log
