# tests/images/valkey/supervisor_senticache.conf
[program:senticache]
command=/usr/bin/valkey-senticache /etc/valkey/senticache.conf
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=TERM
priority=5
redirect_stderr=true
stdout_logfile=/var/log/valkey/senticache.log
