# tests/images/valkey/supervisor_rdsync.conf
[program:rdsync]
command=/usr/bin/rdsync
process_name=%(program_name)s
autostart=true
startretries=100000
autorestart=true
stopsignal=TERM
priority=3
redirect_stderr=true
stdout_logfile=/var/log/rdsync.log
