# tests/images/valkey/supervisor_valkey.conf
[program:valkey]
command=/usr/bin/valkey-server /etc/valkey/valkey.conf
process_name=%(program_name)s
autostart=true
autorestart=false
stopsignal=TERM
priority=5
redirect_stderr=true
stdout_logfile=/var/log/valkey/server.log
